# conf/supervisor.conf
[program:gunicorn]
command=/usr/bin/gunicorn application:app -b localhost:8000
directory=/src/homemonitor
user=nobody
autostart=true
autorestart=true

[program:nginx]
command=/usr/sbin/nginx
autostart=true
autorestart=true
